# root/etc/systemd/system/live_autologin.service
[unit]
Description="Starting Xorg for Live Media"
[Service]
User=lilyum
Execstart=startx .xinitrc
[Install]
WantedBy=multi-user.target

